# root/etc/services.d/flexget/run
#!/usr/bin/with-contenv bash

exec \
      s6-setuidgid abc /usr/local/bin/flexget -c /config/config.yml \
      --logfile ${LOG_FILE} --loglevel ${LOG_LEVEL} daemon start --autoreload-config
